# conf/alquotes.supervisor.conf
[program:alquotes]
command=/home/sandbox/alice-skill-quotes/env/bin/uwsgi --ini /home/sandbox/alice-skill-quotes/conf/alquotes.uwsgi.ini


